# private_dot_config/private_systemd/private_user/mail-sync.timer
[Unit]
Description=Periodically sync mail

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
AccuracySec=1min

[Install]
WantedBy=timers.target
