# private_dot_config/private_systemd/private_user/mail-sync.timer
[Unit]
Description=Periodically sync mail

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
