# private_dot_config/private_systemd/private_user/mail-sync.service
[Unit]
Description=Sync mail via mbsync and index with notmuch
After=network-online.target
Wants=network-online.target
ConditionPathIsDirectory=%h/data/mail

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mail-sync
StandardOutput=journal
StandardError=journal
