# private_dot_config/private_systemd/private_user/rclone-mount@.service.tmpl
[Unit]
Description=rclone FUSE mount for %i
Documentation=https://rclone.org/commands/rclone_mount/
After=network-online.target
Wants=network-online.target
ConditionPathExists=%h/.config/rclone/rclone.conf

[Service]
Type=notify
ExecStartPre=/usr/bin/mkdir -p %h/Cloud/%i
ExecStart={{ lookPath "rclone" | default (printf "%s/bin/rclone" .homebrew.prefix) }} mount %i: %h/Cloud/%i \
    --vfs-cache-mode full \
    --vfs-cache-max-size 5G \
    --vfs-cache-max-age 24h \
    --dir-cache-time 1h \
    --buffer-size 256M \
    --vfs-read-chunk-size 32M \
    --vfs-read-chunk-size-limit 1G \
    --transfers 8 \
    --poll-interval 15m \
    --retries 10 \
    --low-level-retries 20 \
    --retries-sleep 5s \
    --timeout 30m \
    --contimeout 60s
ExecStop=/usr/bin/fusermount3 -uz %h/Cloud/%i
Restart=always
RestartSec=15
RestartSteps=5
RestartMaxDelaySec=5m
StartLimitIntervalSec=10m
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
