# private_dot_config/private_systemd/private_user/rclone-mount@.service.tmpl
[Unit]
Description=rclone FUSE mount for %i
Documentation=https://rclone.org/commands/rclone_mount/
After=network.target
ConditionPathExists=%h/.config/rclone/rclone.conf

[Service]
Type=notify
ExecStartPre=/usr/bin/mkdir -p %h/Cloud/%i
ExecStart={{ lookPath "rclone" | default (printf "%s/bin/rclone" .homebrew.prefix) }} mount %i: %h/Cloud/%i \
    --vfs-cache-mode full \
    --vfs-cache-max-size 5G \
    --vfs-cache-max-age 24h \
    --dir-cache-time 1h \
    --buffer-size 256M \
    --vfs-read-chunk-size 32M \
    --vfs-read-chunk-size-limit 1G \
    --transfers 8 \
    --poll-interval 5m
ExecStop=/usr/bin/fusermount3 -uz %h/Cloud/%i
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
