# private_dot_config/private_systemd/private_user/vdirsyncer.timer
[Unit]
Description=Periodically sync calendars, contacts and tasks

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
