# system-ansible/roles/btrbk/files/btrbk-root.timer
[Unit]
Description=btrbk /root snapshot timer (6 hour interval)

[Timer]
OnBootSec=1min
OnUnitActiveSec=6h
AccuracySec=10min
Persistent=true

[Install]
WantedBy=timers.target
