# private_dot_config/private_systemd/private_user/homebrew.vdirsyncer.service.d/override.conf
[Unit]
ConditionPathExists=%h/.config/vdirsyncer/config

[Service]
StandardOutput=journal
StandardError=journal
