# dot_config/private_systemd/private_user/atuin.service
[Unit]
Description=Atuin daemon

[Service]
Type=simple
ExecStart=%h/.local/bin/atuin daemon
Restart=always

[Install]
WantedBy=default.target
