# dot_config/private_systemd/private_user/msmtp-runqueue.service
[Unit]
Description=Send mails from msmtp queue
[Service]
Type=oneshot
ExecStart=%h/bin/mail/msmtp-runqueue.sh
