# dot_config/private_systemd/private_user/readonly_msmtp-runqueue.service
[Unit]
Description=Send mails from msmtp queue
ConditionPathExists=%h/.local/share/msmtpqueue
[Service]
Type=oneshot
ExecStart=%h/bin/mail/msmtp-runqueue.sh
