# dot_config/private_systemd/private_user/restic-backup.timer
[Unit]
Description=Backup with restic every hour
[Timer]
OnStartupSec=5m
OnUnitActiveSec=30m
Persistent=true
[Install]
WantedBy=timers.target
