# dot_config/containers/systemd/backrest.container
[Unit]
Description=Backrest - Web UI for restic backup
Documentation=https://github.com/garethgeorge/backrest

[Container]
Image=docker.io/garethgeorge/backrest:v1.9.1
AutoUpdate=registry
Label=io.containers.autoupdate=registry
HostName=%H

# Persist backrest data and configuration
Volume=%h/.local/share/backrest:/data:Z
Volume=%h/.config/backrest:/config:Z
Volume=%h/.cache/backrest:/cache:Z

# Mount home directory for backup access (read-only for safety)
# Note: No SELinux relabeling (:Z) to avoid relabeling entire $HOME
Volume=%h:/backup-home:ro

# Network - localhost only for security
PublishPort=127.0.0.1:9898:9898

# Environment
Environment=BACKREST_DATA=/data
Environment=BACKREST_CONFIG=/config/config.json
Environment=XDG_CACHE_HOME=/cache

# Security - container level
# SELinux disabled for this container to allow reading entire $HOME for backups
SecurityLabelDisable=true
ReadOnlyTmpfs=true
NoNewPrivileges=true
DropCapability=ALL

# Resource limits
Memory=2G

[Service]
# Create required directories before starting
ExecStartPre=/usr/bin/mkdir -p %h/.local/share/backrest
ExecStartPre=/usr/bin/mkdir -p %h/.config/backrest
ExecStartPre=/usr/bin/mkdir -p %h/.cache/backrest
CPUQuota=200%
Restart=always
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=default.target
