# dot_config/containers/systemd/syncthing.container
[Unit]
Description=Syncthing - File synchronization
Documentation=https://docs.linuxserver.io/images/docker-syncthing/
After=network-online.target
Wants=network-online.target

[Container]
Image=lscr.io/linuxserver/syncthing:v2.0.10-ls198
AutoUpdate=registry
Label=io.containers.autoupdate=registry

# Environment - LinuxServer.io standard
Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=Etc/UTC

# Persist syncthing config and data
Volume=%h/.config/syncthing:/config:Z

# Mount data directories for synchronization
# Adjust these paths based on what you want to sync
Volume=%h/sync:/sync:Z

# Network ports
# 8384: Web UI
# 22000: TCP file transfers
# 22000/udp: QUIC file transfers
# 21027/udp: Local discovery
PublishPort=127.0.0.1:8384:8384
PublishPort=22000:22000
PublishPort=22000:22000/udp
PublishPort=21027:21027/udp

# Security - container level
User=keep-id
SecurityLabelDisable=false
ReadOnlyTmpfs=true
NoNewPrivileges=true
DropCapability=ALL

# Resource limits
Memory=1G
CPUQuota=200%

[Service]
Restart=always
TimeoutStartSec=300
TimeoutStopSec=30

# Systemd security hardening
# Network - needs external access for sync
PrivateNetwork=false

# Filesystem
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/syncthing
ReadWritePaths=%h/sync
# Add more ReadWritePaths as needed for additional sync folders

# Kernel
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectProc=invisible

# Misc hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectClock=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete

[Install]
WantedBy=default.target
