# private_dot_config/exact_private_autostart/readonly_solaar.desktop
[Desktop Entry]
Type=Application
Name=Solaar
Exec=solaar -w hide
